# site/profile/templates/gunicorn.socket.epp
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=<%= $address %>
BindIPv6Only=both

[Install]

WantedBy=sockets.target
